# cryptobot.service
[Unit]
Description=CryptoValueProtecter

[Service]
ExecStart=python CryptoBot.py

[Install]
WantedBy=multi-user.target